# deploy/termag.service
[Unit]
# Edit User, WorkingDirectory, ExecStart, HOME, and EnvironmentFile to match
# your install. The paths below assume a checkout at /opt/termag owned by the
# `termag` user — substitute your own.
Description=termag workspace manager backend
After=network.target postgresql.service

[Service]
Type=simple
User=termag
WorkingDirectory=/opt/termag/backend
ExecStart=/usr/bin/node /opt/termag/backend/dist/index.js
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=10
Environment=NODE_ENV=production
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/termag"
EnvironmentFile=/opt/termag/backend/.env

[Install]
WantedBy=multi-user.target
